# systemd/.config/systemd/user/elephant.service

[Unit]
Description=Elephant
After=default.target

[Service]
Type=simple
ExecStart=elephant
Restart=always

[Install]
WantedBy=default.target
